# apis/auth/service/log/run
#!/bin/sh

exec setuidgid directorofme \
     multilog -t ./main
